# runit/acme/run
#!/bin/sh
exec 1>&2
echo -e "starting acme.sh ..."
