# runit/acme/run
#!/bin/sh
echo -e "starting acme.sh ..."
exec sleep 300
