# runit/acme/run
#!/bin/sh
echo -e 'Starting acme.sh ...'
