# runit/acme/run
#!/bin/sh
exec 2>&1
echo 'Starting acme.sh ...'
exec chpst -e acme.sh
