# runit/acme/run
#!/bin/sh
exec 2>&1
echo -e 'Starting acme.sh ...'
exec chpst -u root /root/ssl.sh
