# runit/acme/run
#!/bin/sh
echo -e 'Starting acme.sh ...'
exec chpst -e /root/ssl.sh
