# runit/acme/run
#!/bin/sh
exec 2>&1
echo -e  "starting x-ui..."
cd /root && exec chpst test.sh
