# runit/acme/run
#!/bin/sh
echo -e "starting acme.sh ..."
exec crontab -l && sleep 60
