# runit/acme/run
#!/bin/sh
echo -e "Acme.sh script is working..."
exec sleep 300
