# runit/acme/run
#!/bin/sh
exec 2>&1
echo -e "starting acme.sh ..."
exec crontab -l
